# packages/qmp-small-node/files/etc/init.d/qmp_autoconf
#!/bin/sh /etc/rc.common
# Copyright (C) 2011 qmp.cat
START=99
CONTROL_FILE="/qmp_configured"
QMP_DIR="/etc/qmp"

configure() {
	echo "Starting qmp configuration..."

	echo "[Configuring wireless]"
	. $QMP_DIR/qmp_wireless.sh
	qmp_configure_wifi_initial
	qmp_configure_wifi
	wifi

	echo "[Configuring networking]"
	. $QMP_DIR/qmp_functions.sh
	qmp_configure
#	/etc/init.d/network restart	                       
#	/etc/init.d/bmx6 restart	
	touch "$CONTROL_FILE"
	echo "Done, rebooting system to complete configuration"
	reboot
}

start() {
if [ ! -f "$CONTROL_FILE" ]; then
	configure
else
	echo "QMP configured. Remove $CONTROL_FILE to force reconfiguration"
fi
}

stop() {
	echo "Nothing to do"
}

restart() {
	stop
	start
}
