# packages/qmp-small-node/files/etc/init.d/qmp_autoconf
#!/bin/sh /etc/rc.common
#    Copyright (C) 2011 Fundacio Privada per a la Xarxa Oberta, Lliure i Neutral guifi.net
#
#    This program is free software; you can redistribute it and/or modify
#    it under the terms of the GNU General Public License as published by
#    the Free Software Foundation; either version 2 of the License, or
#    (at your option) any later version.
#
#    This program is distributed in the hope that it will be useful,
#    but WITHOUT ANY WARRANTY; without even the implied warranty of
#    MERCHANTABILITY or FITNESS FOR A PARTICULAR PURPOSE.  See the
#    GNU General Public License for more details.
#
#    You should have received a copy of the GNU General Public License along
#    with this program; if not, write to the Free Software Foundation, Inc.,
#    51 Franklin Street, Fifth Floor, Boston, MA 02110-1301 USA.
#
#    The full GNU General Public License is included in this distribution in
#    the file called "COPYING".

START=99
CONTROL_FILE="/qmp_configured"
QMP_DIR="/etc/qmp"

configure() {
	echo "Starting qmp configuration..."
	echo "Starting qmp configuration..." > /dev/console

	echo "[Configuring wireless]"
	. $QMP_DIR/qmp_wireless.sh
	qmp_configure_wifi_initial
	qmp_configure_wifi
	wifi

	echo "[Configuring networking]"
	. $QMP_DIR/qmp_functions.sh
	qmp_configure
	/etc/init.d/network restart	                       
	sleep 5
	QMP_HOSTNAME=$(uci get system.@system[0].hostname)
	echo "Update hostname for: '$QMP_HOSTNAME'"
	echo "$QMP_HOSTNAME" > /proc/sys/kernel/hostname
	/etc/init.d/olsrd restart	
	/etc/init.d/bmx6 restart	
	ifup -a
	touch "$CONTROL_FILE"
#	echo "Done, rebooting system to complete configuration"
#	reboot
}

start() {
if [ ! -f "$CONTROL_FILE" ]; then
	configure
else
	echo "QMP configured. Remove $CONTROL_FILE to force reconfiguration"
fi
}

stop() {
	echo "Nothing to do"
}

restart() {
	stop
	start
}
