# packages/qmp-system/files/etc/uci-defaults/qmp_set_luci_app_bmx6
#!/bin/sh
uci set luci-bmx6.luci=bmx6
uci set luci-bmx6.luci.place="qmp Mesh"
uci set luci-bmx6.luci.position=3
uci set luci-bmx6.luci.json="exec:/www/cgi-bin/bmx6-info -s"
uci set luci-bmx6.luci.ignore=0
uci commit luci-bmx6
