# packages/qmp-system/files/etc/uci-defaults/qmp_remove_old_gateways
#!/bin/sh
i=0
while true; do 
 uci -q del qmp.@gateway[0]
 [ $? -ne 0 ] && break
 i=$(($i+1))
done

uci commit
