# packages/qmp-system/files/etc/uci-defaults/qmp_set_opkg_repo
#!/bin/sh
sed s/12.09.1/12.09-rc2/g -i /etc/opkg.conf
