# packages/qmp-system/files/etc/uci-defaults/qmp_set_luci_app_bmx7
#!/bin/sh
uci set luci-bmx7.luci=bmx7
uci set luci-bmx7.luci.place="qmp Mesh"
uci set luci-bmx7.luci.position=3
uci set luci-bmx7.luci.json="exec:/www/cgi-bin/bmx7-info -s"
uci set luci-bmx7.luci.ignore=0
uci commit luci-bmx7
