# packages/qmp-system/files/etc/uci-defaults/qmp_change_bonus_for_rating
#!/bin/sh
sed -i s/bonus/rating/g /etc/config/gateways 
