# packages/qmp-system/files/etc/uci-defaults/qmp_testing_updates
#!/bin/sh

[ "$(uci get qmp.wireless.country)" == "SP" ] && {
	uci set qmp.wireless.country=US
}

[ "$(uci get qmp.update.url)" == "http://fw.qmp.cat/testing" ] && {
	uci set qmp.update.url="http://fw.qmp.cat/testing"
}

uci -q get qmp.networks.disable_mesh_dhcp || uci set qmp.networks.disable_mesh_dhcp=1
uci commit qmp
