# packages/qmp-system/files/etc/uci-defaults/qmp_banner
#!/bin/sh

cat /etc/qmp/qmp.banner >> /etc/banner
echo " ---------------------------------------------------" >> /etc/banner
echo " $COD ($REL, $BRANCH rev.$REV-$DATE)" >> /etc/banner
echo " ---------------------------------------------------" >> /etc/banner
echo "  * Takeoff clearance: http://en.wikipedia.org/wiki/Air_traffic_control" >> /etc/banner
echo "" >> /etc/banner
echo " http://qmp.cat/projects/qmp/repository/show?branch=$BRANCH&rev=$REV" >> /etc/banner
echo " ---------------------------------------------------" >> /etc/banner

