# packages/qmp-system/files/etc/uci-defaults/qmp_add_mini_snmpd_to_preserve
#!/bin/sh
mini_snmpd="/etc/config/mini_snmpd"
preserve="$(uci get qmp.update.preserve)"
found=0
for p in $preserve; do
        [ "$p" == "$mini_snmpd" ] && { found=1; break; }
done

[ $found -eq 0 ] && {
        uci set qmp.update.preserve="$preserve $mini_snmpd"
        uci commit qmp
}
