# packages/qmp-system/files/etc/uci-defaults/qmp_set_mini_snmpd
#!/bin/sh

[ -f /etc/config/mini_snmpd ] && {
        uci set mini_snmpd.@mini_snmpd[0].enabled=1
        if ! uci get mini_snmpd.@mini_snmpd[0].interfaces | egrep -q 'br-lan|wlan0'        
                then
                uci add_list mini_snmpd.@mini_snmpd[0].interfaces="wlan0"
                uci add_list mini_snmpd.@mini_snmpd[0].interfaces="br-lan"
        fi
        uci commit

	/etc/init.d/mini_snmpd restart
}

