# packages/qmp-system/files/etc/uci-defaults/qmp_remove_olsr_config_fix_vlan
#!/bin/sh

if uci get qmp.networks.mesh_vid_offset 2>/dev/null; then
	echo "Removing deprecated OLSR configuration and adapting vlan VID configuration"
	uci set qmp.networks.mesh_protocol_vids='bmx6:12'
	uci del qmp.networks.mesh_vid_offset 
	uci del qmp.networks.olsr6_mesh_prefix48 
	uci del qmp.networks.olsr6_ripe_prefix48 
	uci commit qmp
fi

