# packages/qmp-system/files/etc/uci-defaults/qmp_change_altermesh_for_libremap
#!/bin/sh
sed -i s/altermap/libremap/g /etc/config/qmp
