# packages/qmp-system/files/etc/uci-defaults/qmp_add_birth_to_preserve
#!/bin/sh
birth="/etc/qmp/.birth"
preserve="$(uci get qmp.update.preserve)"
found=0
for p in $preserve; do
	[ "$p" == "$birth" ] && { found=1; break; }
done

[ $found -eq 0 ] && { 
	uci set qmp.update.preserve="$preserve $birth"
	uci commit qmp
}
