# packages/qmp-system/files/etc/uci-defaults/qmp_b6m_set_default_spread_interval
#!/bin/sh

if ! uci get b6m-spread.general.interval 2>/dev/null; then
	uci set b6m-spread.general.interval=60
fi

