# packages/qmp-system/files/etc/uci-defaults/qmp_000_copy_config_if_not_exist
#!/bin/sh
[ ! -f "/etc/config/qmp" ] && cp -f /etc/qmp/qmp.default-config /etc/config/qmp
