# packages/qmp-system/files/etc/uci-defaults/qmpfix_20120921
#!/bin/sh
uci delete qmp.networks.bmx6_mesh_prefix48
uci delete bmx6.general.globalPrefix
uci set qmp.networks.bmx6_tablePrefTuns=7000
uci set qmp.networks.bmx6_tableTuns=70
uci set qmp.networks.bmx6_throwRules=1
uci commit
