# packages/qmp-system/files/etc/uci-defaults/qmpfix_20130131
#!/bin/sh
# Removes the wifi commnad introduced in rc.local for compatiblity with backfire
sed s/^wifi$//g -i /etc/rc.local
