# files/etc/init.d/qmp-fix-vlan
#!/bin/sh /etc/rc.common
# Copyright (C) 2006 OpenWrt.org

START=20
STOP=98

start() {
	[ -f /proc/net/vlan/config ] && vconfig set_name_type DEV_PLUS_VID_NO_PAD
}

stop() {
	echo
	
}

