# packages/synctincvpn/files/etc/init.d/synctincvpn
#!/bin/sh /etc/rc.common
#
# Copyright © 2012 Routek S.L.
#
# This program is free software; you can redistribute it and/or modify
# it under the terms of the GNU General Public License as published by
# the Free Software Foundation; either version 2 of the License, or
# (at your option) any later version.
#
# This program is distributed in the hope that it will be useful,
# but WITHOUT ANY WARRANTY; without even the implied warranty of
# MERCHANTABILITY or FITNESS FOR A PARTICULAR PURPOSE.  See the
# GNU General Public License for more details.
#
# You should have received a copy of the GNU General Public License along
# with this program. If not, see <http://www.gnu.org/licenses/>.
#
# Authors:
#	Agustí Moll i Garcia
#	Simó Albert i Beltran
#

START=97
STOP=10

start() {

	if ! cat /etc/crontabs/root | grep stvclient > /dev/null; then
		echo "0 * * * * /usr/sbin/stvclient >> /tmp/log/stvclient.log" >> /etc/crontabs/root
		/etc/init.d/cron enable
		/etc/init.d/cron start
	fi
	if [ ! -f /etc/config/tincvpn ]; then
		echo "Generate Config File [/etc/config/tincvpn]."
		generate_configfile
	fi
	/usr/sbin/stvclient
	}

stop() {
	sed -i -e '/^\([^\t ]\+[\t ]\+\)\{5\}\/usr\/sbin\/stvclient >> \/tmp\/log\/stvclient\.log$/d' /etc/crontabs/root
	/etc/init.d/cron restart
	}

reload() {
	start
	}

generate_configfile() {
	cat > /etc/config/tincvpn <<EOF
config tincvpn
	option enable 	1
EOF
	uci batch <<-EOF
	set tincvpn.@tincvpn[0].name=qmptest
	set tincvpn.@tincvpn[0].url=http://tincvpn.qmp.cat:8001/nodes/get
	set tincvpn.@tincvpn[0].key=mykey
        commit tincvpn
EOF

}
