# packages/qmp-system/files/etc/init.d/qmp
#!/bin/sh /etc/rc.common
#
# Copyright © 2016
# Fundació Privada per a la Xarxa Oberta, Lliure i Neutral Guifi.net

START=99

QMP_CONFIG="/etc/config/qmp"

start() {

	if [ ! -e "$QMP_CONFIG" ]; then
		initialize
	fi
}

stop() {

}

restart() {
	start
}

initialize() {
	qmpcontrol initialize
}