# packages/qmp-small-node/files/etc/uci-defaults/qmp_splash_fw
#!/bin/sh

cat << QMP >> /etc/firewall.user
iptables -t nat -A POSTROUTING -s 10.0.0.0/8 -d 10.0.0.0/8 -j RETURN
iptables -t nat -A POSTROUTING -j MASQUERADE
QMP
