# packages/qmp-small-node/files/etc/uci-defaults/qmp_clamp_tcpmss
#!/bin/sh

if cat /etc/firewall.user | grep -e "^# Clamp MSS TCP rule to fix MTU problems"; then exit 0; fi

cat << EOF >> /etc/firewall.user
# Clamp MSS TCP rule to fix MTU problems
iptables -A FORWARD -p tcp -o bmx6_+ -m tcp --tcp-flags SYN,RST SYN -j TCPMSS --clamp-mss-to-pmtu
EOF
