# packages/qmp-small-node/files/etc/uci-defaults/qmp_bmx6_health_patch
#!/bin/sh

if ! cat /etc/crontabs/root | grep bmx6health; then
	echo "* * * * * /etc/qmp/bmx6health.sh >> /tmp/log/bmx6health.log" >> /etc/crontabs/root
	/etc/init.d/cron enable
	/etc/init.d/cron start
fi
