# packages/qmp-small-node/files/etc/uci-defaults/qmp_splash_fw
#!/bin/sh

echo "/usr/sbin/iptables -t nat -N captive-portal" >> /etc/firewall.user
echo "/usr/sbin/iptables -t nat -A PREROUTING -i br-lan -j captive-portal" >> /etc/firewall.user
echo "/usr/sbin/iptables -t nat -A captive-portal -p tcp --destination-port 80 -j DNAT --to-destination $(uci get network.lan.ipaddr):8888" >> /etc/firewall.user

