# packages/qmp-small-node/files/etc/uci-defaults/qmp_banner
#!/bin/sh

cat /etc/qmp/qmp.banner >> /etc/banner
echo "Versio ratafia 1.0 (Sat Oct 29 2011) ---------------" >> /etc/banner
echo "Ratafia Recipe: http://en.wikipedia.org/wiki/Ratafia" >> /etc/banner
echo "----------------------------------------------------" >> /etc/banner


