# packages/qmp-small-node/files/etc/uci-defaults/qmp_banner
#!/bin/sh

cat /etc/qmp/qmp.banner >> /etc/banner
echo "Versio XXXX" >> /etc/banner
echo "----------------------------------------------------" >> /etc/banner


