# packages/qmp-small-node/files/etc/uci-defaults/qmp_disable_b6m
#!/bin/sh

echo "Disabling b6m..."
/etc/init.d/b6m-spread stop
/etc/init.d/b6m-spread disable

