# packages/qmp-small-node/files/etc/uci-defaults/qmp_splash_fw
#!/bin/sh

cat << QMP >> /etc/firewall.user
/usr/sbin/iptables -t nat -A POSTROUTING -j MASQUERADE
QMP
