# packages/qmp-small-node/files/etc/uci-defaults/qmp_disable_fw
#!/bin/sh

uci delete firewall.@zone[1]
uci delete firewall.@zone[0]
uci delete firewall.@defaults[0]
uci delete firewall.@forwarding[0]
uci delete firewall.@rule[1]      
uci delete firewall.@rule[1]
uci commit

