# packages/qmp-small-node/files/etc/uci-defaults/qmpfix_20120921

uci delete bmx6.general.globalPrefix
uci commit
