# packages/luci-theme-qmp/files/etc/uci-defaults/luci-theme-qmp
#!/bin/sh
uci batch <<-EOF
	set luci.themes.qMp=/luci-static/qmp
        commit luci
EOF
	
